# docker-utils/supervisord/php7.1-fpm.conf
[program:php7.1-fpm]
command=/usr/sbin/php7.1-fpm -F
stdout_events_enabled=true
stderr_events_enabled=true
autostart=true
autorestart=true

